# deploy/systemd/mtpanel.service
[Unit]
Description=MTPanel - TeleMT Management Panel
Documentation=https://github.com/NikitaKHS/mtpanel
After=network-online.target
Wants=network-online.target
PartOf=telemt.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/mtpanel
EnvironmentFile=-/etc/mtpanel/mtpanel.env
ExecStart=/opt/mtpanel/mtpanel
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=30s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mtpanel

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/mtpanel /etc/mtpanel /etc/telemt /opt/telemt /etc/systemd/system
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictNamespaces=true
SystemCallFilter=@system-service
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

CapabilityBoundingSet=
AmbientCapabilities=

ProtectProc=invisible
ProcSubset=pid

[Install]
WantedBy=multi-user.target
