# deploy/systemd/mtpanel.service
[Unit]
Description=MTPanel
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/mtpanel
Environment=MTPANEL_CONFIG=/etc/mtpanel/config.json
ExecStart=/opt/mtpanel/mtpanel
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
