# deploy/systemd/mtpanel.service
[Unit]
Description=MTPanel - TeleMT Management Panel
Documentation=https://github.com/NikitaKHS/mtpanel
After=network-online.target
Wants=network-online.target
PartOf=telemt.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/mtpanel
EnvironmentFile=-/etc/mtpanel/mtpanel.env
ExecStart=/opt/mtpanel/mtpanel
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=30s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mtpanel

NoNewPrivileges=false
PrivateTmp=false
PrivateDevices=false
ProtectSystem=off
ProtectHome=off
ReadWritePaths=
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectKernelLogs=false
ProtectControlGroups=false
RestrictSUIDSGID=false
RestrictRealtime=false
LockPersonality=false
MemoryDenyWriteExecute=false
RestrictNamespaces=false
SystemCallFilter=
SystemCallArchitectures=
SystemCallErrorNumber=

CapabilityBoundingSet=
AmbientCapabilities=

ProtectProc=default
ProcSubset=all

[Install]
WantedBy=multi-user.target
